# copy_files/cryptswap1.conf
start on started mountall
script
	/sbin/cryptdisks_start cryptswap1
	/sbin/swapon /dev/mapper/cryptswap1
end script
